# upstart/default-clinicaldata-upstart.conf
# Default ClinicalData upstart config

description "Default Mediator"

# logs to /var/log/upstart/default-clinicaldata.log
console log

start on runlevel [2345]
stop on runlevel [!2345]

respawn

setuid root
setgid root

# Change this to the proper directory of Default-ClinicalData
env CLINICAL_DATA=/home/MEDIATOR-CLINICAL-DATA/Default-ClinicalData

script
  cd $CLINICAL_DATA
  exec bash -c 'source env/bin/activate && gunicorn --worker-class gthread -w 1 --threads 1 app:app'
end script